# test-rktnetes/crash.service
[Service]
ExecStart=/opt/test-rktnetes/crashloop.sh
Restart=always
