# test-rktnetes/test.service
[Service]
ExecStart=/opt/test-rktnetes/launch_kill.sh
Restart=always
